# config_files/80-ax150.rules
ACTION=="add", SUBSYSTEM=="pci", ATTR{device}=="0x2723", RUN+="/usr/local/bin/pcie_set_speed.sh 03:00.0 2"
